# KopiaBackup/Build/test.service
[Unit]
Description=MyTool Backup Dienst
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/mytool backup /home/user/Dokumente
Restart=on-failure
User=myuser

[Install]
WantedBy=multi-user.target
